# system/cron.d/rsnapshot
#
# ts4z
# hook for starting rsnapshot.
#

0 * * * *         root /usr/bin/rsnapshot hourly

50 23 * * *         root /usr/bin/rsnapshot daily

40 23 * * 6         root /usr/bin/rsnapshot weekly

30 23 1 * *         root /usr/bin/rsnapshot monthly

