# system/cron.d/tarsnap
20 10 * * *     root  /local/etc/tarsnap-backup.sh
